# docs/testing/mt5-vps-journal/deploy/systemd/mt5-headless.service
[Unit]
Description=Headless MT5 terminal under Wine/Xvfb
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mt5journal
Group=mt5journal
EnvironmentFile=/etc/mt5-journal.env
WorkingDirectory=/opt/mt5-vps-journal
ExecStart=/opt/mt5-vps-journal/scripts/headless_mt5_service.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
