# docs/testing/mt5-vps-journal/deploy/systemd/mt5-export-sync.timer
[Unit]
Description=Run MT5 export CSV sync every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=mt5-export-sync.service

[Install]
WantedBy=timers.target
